# systemd/folderr-reset-user-resource-usage.timer
[Unit]
Description=Timer for resetting user resource usage

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
