# systemd/folderr-repeat-tasks.service
[Unit]
Description=Reset user resource usage

[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/usr/bin/folderr-repeat-tasks.sh
