# systemd/folderr-create-task-reminders.timer
[Unit]
Description=Timer for creating task reminders

[Timer]
OnBootSec=6h
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
